# scripts/morning-planner.service
[Unit]
Description=Claude Code morning planner

[Service]
Type=oneshot
User=ec2-user
Environment=HOME=/home/ec2-user
Environment=PATH=/home/ec2-user/.local/bin:/home/ec2-user/.nvm/versions/node/current/bin:/usr/local/bin:/usr/bin:/bin
WorkingDirectory=/home/ec2-user
ExecStart=/bin/bash /home/ec2-user/llm-skills/scripts/run-morning-planner.sh
TimeoutStartSec=600
