# scripts/nightly-sync.service
[Unit]
Description=Claude Code nightly project sync

[Service]
Type=oneshot
User=ec2-user
Environment=HOME=/home/ec2-user
Environment=PATH=/home/ec2-user/.local/bin:/home/ec2-user/.nvm/versions/node/current/bin:/usr/local/bin:/usr/bin:/bin
WorkingDirectory=/home/ec2-user
ExecStart=/bin/bash /home/ec2-user/llm-skills/scripts/run-nightly-sync.sh
TimeoutStartSec=600
